# systemd.service
# Systemd service file for Williams Performance Friction
# Place this file at: /etc/systemd/system/williams-performance-friction.service

[Unit]
Description=Williams Performance Friction Flask App
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/app  # Replace with your actual app path
Environment=PATH=/path/to/your/app/venv/bin  # Replace with your venv path
Environment=DATABASE_URL=postgresql://username:password@localhost/williams_db  # Configure your database
Environment=SESSION_SECRET=your-super-secret-session-key-here  # Generate a strong secret
Environment=SENDGRID_API_KEY=your-sendgrid-api-key  # Your SendGrid API key
ExecStart=/path/to/your/app/venv/bin/gunicorn --config gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target